# querio.service
[Unit]
Description=Querio RAG API Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ubuntu/querio
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
User=ubuntu
Group=ubuntu
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
